# packaging/debian/DEBIAN/prerm
#!/bin/sh

PACKAGE="proactive-agent"

# Service removal
/etc/init.d/$PACKAGE stop
\update-rc.d -f $PACKAGE remove
