# packaging/rpm/prerm
#!/bin/sh

PACKAGE="proactive-agent"

/sbin/service $PACKAGE stop || true
/sbin/chkconfig --del $PACKAGE
