# packaging/deb/prerm
#!/bin/sh

PACKAGE="proactive-agent"

invoke-rc.d $PACKAGE stop
update-rc.d -f $PACKAGE remove
